# noauto_mounter.service
# /etc/systemd/system/auto_share.service
[Unit]
Description=Check fstab noauto mounts

[Service]
Type=simple
ExecStart=/usr/bin/noauto_mounter
